# deploy/vm-health-watchdog/vm-health-watchdog.timer
[Unit]
Description=Run vm-health-watchdog every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Unit=vm-health-watchdog.service

[Install]
WantedBy=timers.target
